# trunk/framework/aspdotnet/ki/noninteractive/.crontab
#
# $Id$
#
00 00 * * * wget --output-document=/dev/null https://sslws010.alentus.com/frompaper2web/KiAspdotnetFramework/noninteractive/milestone_sweep.aspx
